# 10-pi-zero-2-met-mqtt-broker/etc/systemd/system/mqtt-status.service
[Unit]
After=mosquitto.service

[Service]
User=ubuntu
Group=ubuntu
ExecStart=/home/ubuntu/mqtt/status.sh

[Install]
WantedBy=default.target
